# build/gdelt-live.service
[Unit]
Description=Run gdelt-live.py

[Service]
Type=oneshot
ExecStart=python /opt/gdelt-diff/core/gdelt-live.py
