# build/gdelt-live.timer
[Unit]
Description=Run gdelt-live.service

[Timer]
Unit=gdelt-live.service
OnBootSec=5min
OnUnitActiveSec=1min

[Install]
WantedBy=multi-user.target
